# src/rne/systemd/rne-worker.service
[Unit]
Description=rne encoder worker
After=network.target
ConditionPathExists=/mnt/media

[Service]
Type=simple
User=rip
Group=mediagroup
ExecStart=__RNE_WORKER_BIN__
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
Nice=10
Environment=HOME=/home/rip
Environment=RNE_DB=/home/rip/.local/state/rne/jobs.db

[Install]
WantedBy=multi-user.target
