# src/rne/systemd/rne-dashboard.service
[Unit]
Description=rne dashboard
After=network.target
ConditionPathExists=/mnt/media

[Service]
Type=simple
ExecStart=__RNE_DASHBOARD_BIN__
Restart=on-failure
RestartSec=5
Environment=HOME=/home/rip
Environment=RNE_DB=/home/rip/.local/state/rne/jobs.db

[Install]
WantedBy=default.target
