# src/rne/systemd/rne-worker.service
[Unit]
Description=rne encoder worker
After=network.target
ConditionPathExists=/mnt/media

[Service]
Type=simple
ExecStart=__RNE_WORKER_BIN__
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
Nice=10

[Install]
WantedBy=default.target
